# scripts/com.stitchflow.niimbot-agent.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.stitchflow.niimbot-agent</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>/Users/your-user/Downloads/stitchflow (1)/scripts/niimbot-agent.mjs</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>NIIMBOT_AGENT_TOKEN</key>
      <string>replace-with-strong-token</string>
      <key>NIIMBOT_AGENT_PORT</key>
      <string>9130</string>
      <key>NIIMBOT_ALLOWED_ORIGINS</key>
      <string>http://localhost:5173,http://127.0.0.1:5173</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/stitchflow-niimbot-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/stitchflow-niimbot-agent.err.log</string>
  </dict>
</plist>
